# ssr.conf
[program:ssr]
command=python /opt/shadowsocksr/server.py
autorestart=true
autostart=true
user=root
